# config/launchd/com.lifeos.crm-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lifeos.crm-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nathanramia/Documents/Code/LifeOS/.venv/bin/python</string>
        <string>/Users/nathanramia/Documents/Code/LifeOS/scripts/run_all_syncs.py</string>
        <string>--force</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nathanramia/Documents/Code/LifeOS</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/nathanramia</string>
        <key>PATH</key>
        <string>/Users/nathanramia/Documents/Code/LifeOS/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/nathanramia/Documents/Code/LifeOS</string>
    </dict>

    <!-- Run in Aqua session for better permissions (iMessage, Contacts access) -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Run daily at 6 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Also run on load (after reboot) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Log files -->
    <key>StandardOutPath</key>
    <string>/Users/nathanramia/Documents/Code/LifeOS/logs/crm-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nathanramia/Documents/Code/LifeOS/logs/crm-sync-error.log</string>

    <!-- Don't restart on failure (will retry next day) -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
